# src/main/resources/run/flume-supervisor.conf
[program:flume]
command=sh /opt/apache-flume/bin/start-flume.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/stdout.log
autostart=true
autorestart=true